# wpa_supplicant.service.d/override.conf
[Service]
CapabilityBoundingSet=~CAP_LINUX_IMMUTABLE
ProtectHome=yes
LockPersonality=yes
SystemCallFilter=~@cpu-emulation @obsolete
SystemCallArchitectures=native
