# power-profiles-daemon.service.d/override.conf
[Service]
CapabilityBoundingSet=~CAP_LINUX_IMMUTABLE CAP_NET_ADMIN
ProtectHostname=yes
PrivateUsers=yes
SystemCallFilter=~@cpu-emulation @obsolete
SystemCallArchitectures=native
IPAddressDeny=any
