# NetworkManager.service.d/override.conf
[Service]
SystemCallFilter=~@cpu-emulation @debug @obsolete
SystemCallArchitectures=native
