# realmd.service.d/override.conf
[Service]
CapabilityBoundingSet=~CAP_LINUX_IMMUTABLE
ProtectHome=yes
LockPersonality=yes
SystemCallFilter=~@cpu-emulation @debug @mount @obsolete @swap
SystemCallArchitectures=native
