# plymouth-start.service.d/override.conf
[Service]
CapabilityBoundingSet=~CAP_LINUX_IMMUTABLE
ProtectHostname=yes
ProtectClock=yes
LockPersonality=yes
SystemCallFilter=~@cpu-emulation @obsolete
SystemCallArchitectures=native
