# low-memory-monitor.service.d/override.conf
[Service]
CapabilityBoundingSet=~CAP_LINUX_IMMUTABLE
ProtectHostname=yes
SystemCallFilter=~@obsolete
