# NetworkManager-wait-online.service.d/override.conf
[Service]
CapabilityBoundingSet=~CAP_LINUX_IMMUTABLE
ProtectHome=yes
SystemCallFilter=~@cpu-emulation @mount @obsolete @swap
SystemCallArchitectures=native
