# uresourced.service.d/override.conf
[Service]
CapabilityBoundingSet=~CAP_LINUX_IMMUTABLE
ProtectHostname=yes
SystemCallArchitectures=native
IPAddressDeny=any
