# rtkit-daemon.service.d/override.conf
[Service]
ProtectHostname=yes
SystemCallFilter=~@cpu-emulation @obsolete
SystemCallArchitectures=native
