# udisks2.service.d/override.conf
[Service]
SystemCallArchitectures=native
IPAddressDeny=any
